# deployment/supervisor/microblog.conf
[program:DMS]
command=/home/ubuntu/DMS/venv/bin/gunicorn -b localhost:8000 -w 4 DMS:app
directory=/home/ubuntu/DMS
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
