# deployment/supervisor/microblog-tasks.conf
[program:DMS-tasks]
command=/home/ubuntu/DMS/venv/bin/rq worker DMS-tasks
numprocs=1
directory=/home/ubuntu/DMS
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
